# systemd/ad-password-manager.service
[Unit]
Description=Helixon Password Self-Service Reset
After=network.target postgresql.service

[Service]
Type=simple
User=www-data
Group=www-data
# ↓↓↓ 部署时把下面两行的路径改为实际项目路径 ↓↓↓
WorkingDirectory=/opt/helixon-pwd/backend
Environment="PATH=/opt/helixon-pwd/backend/.venv/bin"
ExecStart=/opt/helixon-pwd/backend/.venv/bin/gunicorn \
          --workers 4 --threads 2 \
          --bind 127.0.0.1:5000 \
          app:app
Restart=always
RestartSec=10

# 日志
StandardOutput=append:/opt/helixon-pwd/logs/gunicorn.log
StandardError=append:/opt/helixon-pwd/logs/gunicorn_error.log

# 安全
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
