# systemd/ad-password-manager.service
[Unit]
Description=AD Password Management System
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/path/to/your/project/backend
Environment="PATH=/path/to/your/project/backend/venv/bin"
ExecStart=/path/to/your/project/backend/venv/bin/gunicorn --bind 0.0.0.0:5000 --workers 4 --threads 2 app:app
Restart=always
RestartSec=10

# 日志
StandardOutput=append:/path/to/your/project/logs/gunicorn.log
StandardError=append:/path/to/your/project/logs/gunicorn_error.log

# 安全设置
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
